# services/nginx/run
#!/usr/bin/env bash

exec /sbin/setuser www-data nginx
